# debian/prerm
#/usr/bin/env sh

systemctl disable sj201.service
